# .config/systemd/user/udiskie.service
[Unit]
Description=udiskie automounter
Documentation=man:udiskie(8)

[Service]
Type=simple
ExecStart=/usr/bin/udiskie
